# lb/supervisord.conf
[supervisord]
nodaemon=true
loglevel=critical

[program:sshd]
command=/usr/sbin/sshd -D
autostart=true
autorestart=true
stderr_logfile=/var/log/sshd.err.log
stdout_logfile=/var/log/sshd.out.log

[program:nginx]
command=/usr/sbin/nginx -g "daemon off;"
autostart=true
autorestart=true
stderr_logfile=/var/log/nginx.err.log
stdout_logfile=/var/log/nginx.out.log

[program:corosync]
command=/usr/sbin/corosync
autostart=true
autorestart=true
stderr_logfile=/var/log/corosync.err.log
stdout_logfile=/var/log/corosync.out.log

[program:pacemaker]
command=/usr/sbin/pacemakerd
autostart=true
autorestart=true
stderr_logfile=/var/log/pacemaker.err.log
stdout_logfile=/var/log/pacemaker.out.log

[program:pcsd]
command=/usr/sbin/pcsd
autostart=true
autorestart=true
stderr_logfile=/var/log/pcsd.err.log
stdout_logfile=/var/log/pcsd.out.log
